# build/overlay/etc/periodic/daily/updateblocklists
#!/bin/bash
sleep $(( RANDOM % 500 ))
/opt/scripts/update_blocklists.sh
